# bastion/ansible/roles/iamsync/templates/iamsync.timer.j2
[Unit]
Description=Iamsync timer

[Timer]
OnCalendar={{ iamsync_cron }}
Persistent=true

[Install]
WantedBy=timers.target
